# config/launchd/com.eddieb.dev-health-morning-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eddieb.dev-health-morning-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eddiebelaval/Development/generate-morning-report.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/eddiebelaval/Development/reports/launchd-morning-report.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eddiebelaval/Development/reports/launchd-morning-report-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/Users/eddiebelaval/.nvm/versions/node/v20.18.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/eddiebelaval</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
